# contrib/systemd/glossa.service
[Unit]
Description=glossa user-session daemon
PartOf=graphical-session.target
After=graphical-session.target dotool.service
Wants=dotool.service

[Service]
Type=simple
EnvironmentFile=-%h/.config/glossa/glossa.env
ExecStart=%h/.local/bin/glossa --config %h/.config/glossa/config.toml daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
